# jassbot.service.d/override.conf
[Service]
Environment="PASSWORD=uX7PixTaUNkE"
